# config/com.justinpaulson.status-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.justinpaulson.status-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-N</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=accept-new</string>
        <string>-R</string>
        <string>0.0.0.0:9998:localhost:9999</string>
        <string>-i</string>
        <string>/Users/justinpaulson/.ssh/corre-prod.pem</string>
        <string>ubuntu@54.221.142.206</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/justinpaulson/Developer/status/log/tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/justinpaulson/Developer/status/log/tunnel.log</string>
</dict>
</plist>
